# ansible/roles/benchmark/templates/benchworker.service
[Unit]
Description=ISUCON9 Benchworker
Requires=docker.service

[Service]
EnvironmentFile=-/home/admin/isucon/.env
WorkingDirectory=/home/admin/isucon/bin

Type=simple
ExecStart=/home/admin/isucon/bin/benchworker run

Restart=always
RestartSec=10s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
