# ansible/roles/benchmark/templates/payment.service
[Unit]
Description=ISUCON9 Payment API
Requires=docker.service

[Service]
EnvironmentFile=-/home/admin/isucon/.env
WorkingDirectory=/home/admin/isucon/bin

Type=simple
ExecStart=/home/admin/isucon/bin/payment
Nice=-1

Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
